# docs/server/gestionale-frontend.service
[Unit]
Description=Gestionale Frontend Next.js
After=network.target

[Service]
Type=simple
User=mauri
WorkingDirectory=/home/mauri/gestionale-fullstack/frontend
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 